# buildroot/buildroot-1.0/ch_conf/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
# 2006-07-26 mf: Initializes the leds for startup procedure
null::sysinit:/root/bin/startled.sh &
#null::sysinit:/bin/mount -o remount,ro /
null::sysinit:/bin/mount -t proc proc /proc
null::sysinit:/bin/mount -a
null::sysinit:/bin/mount -o remount,ro /
null::sysinit:/bin/hostname -F /etc/hostname
null::sysinit:/sbin/ifconfig lo 127.0.0.1 up
null::sysinit:/sbin/route add -net 127.0.0.0 netmask 255.0.0.0 lo
# 2006-05-02 gc set ip
null::sysinit:/sbin/ifconfig eth0 192.168.42.2
# now run any rc scripts
::sysinit:/etc/init.d/rcS
# 2006-07-26 mf: Terminates led startup procedure
null::sysinit:/root/bin/termled.sh &
# Set up a couple of getty's
#tty1::respawn:/sbin/getty 38400 tty1
#tty2::respawn:/sbin/getty 38400 tty2

# Put a getty on the serial port
ttyS2::respawn:/sbin/getty -L ttyS2 38400 ansi-ht
#::respawn:/sbin/getty -L /dev/bluetooth/rfcomm/0 38400 xterm
# hier kein TTY angeben, dann wird automatisch /dev/console genommen!!!
# wieder raus, console ist kein tty und hat kein job control
#::respawn:/sbin/getty -L console 38400 xterm
# 2006-06-08 mf start pppd
# null::respawn:/usr/sbin/pppd call gprs nodetach
# 2006-09-22 mf start mgetty on ttyS0
# ttyS0::respawn:/sbin/mgetty ttyS0
# Logging junk not needed cause started by init.d
#null::sysinit:/bin/touch /var/log/messages
#null::respawn:/sbin/syslogd -s 200kb -b 1 -n -m 0
#null::respawn:/sbin/klogd -n
#tty3::respawn:/usr/bin/tail -f /var/log/messages

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
# 2006-07-26 mf: unmounting SRAM as logging buffer for system reboot
#null::shutdown:/etc/init.d/S11srammount stop
null::shutdown:/usr/bin/killall klogd
null::shutdown:/usr/bin/killall syslogd
null::shutdown:/bin/umount -a -r
#null::shutdown:/sbin/swapoff -a

